# systemd/flask-api.service
[Unit]
Description=Flask API for Blog CMS
After=network.target

[Service]
Type=exec
User=blogcms
Group=blogcms
WorkingDirectory=/opt/blog-cms/flask-api
Environment=PATH=/opt/blog-cms/flask-api/venv/bin
ExecStart=/opt/blog-cms/flask-api/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 3 api:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target